# poisson/rootfs/etc/services.d/tor/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Only start Tor if enabled in config
if bashio::config.true "enable_tor"; then
    # Ensure Tor data directory exists and is writable
    mkdir -p /var/lib/tor
    chmod 700 /var/lib/tor

    # Start with base SOCKS-only config
    cp /etc/torrc /tmp/torrc

    bashio::log.info "Starting Tor daemon..."
    exec tor -f /tmp/torrc
else
    bashio::log.info "Tor is disabled in configuration. Sleeping."
    exec sleep infinity
fi
